# sys/lib/systemd/system/prometheus-netscanner.timer
[Unit]
Description=Scan the network every minute for new hosts

[Timer]
OnStartupSec=60s
OnUnitInactiveSec=60s

[Install]
WantedBy=timers.target
